# backup/backup.cron.sample
# Backup SNCRS data daily.
0 12 * * * <BACKUP_SCRIPT_LOCATION>
